# Apps/UnboundRedis/timers/unbound-rpz.service
[Unit]
Description=Atualização da Lista RPZ
After=network.target

[Service]
ExecStart=/opt/unbound/scripts/update-rpz.sh
