# Apps/UnboundRedis/timers/unbound-blocklist.service
[Unit]
Description=Atualização da Lista de Bloqueios
After=network.target

[Service]
ExecStart=/opt/unbound/scripts/update-blocklists.sh
