# Apps/UnboundRedis/timers/unbound-roothints.service
[Unit]
Description=Atualização dos Root Hints
After=network.target

[Service]
ExecStart=/opt/unbound/scripts/update-roothints.sh
